# cron/crontab
0 0 * * * root wget http://d9/api/every_day?key=YmFubGlzdC5pbmZv